# com.marrow.heart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marrow.heart</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/marrow-core/.venv/bin/marrow</string>
    <string>run</string>
    <string>--config</string>
    <string>/opt/marrow-core/marrow.toml</string>
    <string>--json-logs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/opt/marrow-core</string>

  <key>UserName</key>
  <string>marrow</string>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/marrow/runtime/logs/heart.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/marrow/runtime/logs/heart.stderr.log</string>
</dict>
</plist>
