# com.marrow.heart.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marrow.heart.sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/marrow-core/sync.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>5</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/marrow/runtime/logs/sync.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/marrow/runtime/logs/sync.stderr.log</string>
</dict>
</plist>
