# marrow-heart-sync.service
[Unit]
Description=Marrow sync job
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/marrow-core
Environment=PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
ExecStart=/opt/marrow-core/sync.sh
StandardOutput=append:/Users/marrow/runtime/logs/sync.stdout.log
StandardError=append:/Users/marrow/runtime/logs/sync.stderr.log
